# patches/etc/systemd/system/apprendys.service
[Unit]
Description=Apprendys - Environnement educatif accessible
After=apprendys-mount.service multi-user.target
# Lightdm attend que ce service soit READY (P5 monte, home pret, touchpad rebind)
# Sans Before= : lightdm et apprendys demarraient en parallele -> race condition P5
Before=lightdm.service

[Service]
Type=notify
TimeoutStartSec=30
ExecStart=/opt/apprendys/boot/apprendys-boot.sh
ExecStop=/bin/bash -c "pkill -f apprendys-update.sh 2>/dev/null; kill $(cat /run/gnuramage.pid 2>/dev/null) 2>/dev/null; true"
TimeoutStopSec=15
Restart=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
